# install/macosx/org.cryptable.eidonkey.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
	<dict>
		<key>Label</key>
		<string>org.cryptable.eidonkey</string>
		<key>Program</key>
		<string>/Users/david/Applications/eidonkey/eidonkey</string>
		<key>WorkingDirectory</key>
		<string>/Users/david/Applications/eidonkey/</string>
		<key>RunAtLoad</key>
		<true/>
	</dict>
</plist>